# package/etc/systemd/system/mauircd.service
[Unit]
Description=mauIRC backend

[Service]
ExecStart=/usr/bin/mauircd

[Install]
WantedBy=network.target
