# package/etc/systemd/system/mauirc.service
[Unit]
Description=mauIRC backend

[Service]
ExecStart=/usr/bin/mauirc-server

[Install]
WantedBy=network.target
